# systemd/claude-heartbeat.service
[Unit]
Description=ClaudeHeartbeat - keep the Claude 5-hour window anchored
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Method B (server-truth, recommended). For the simple interval-based
# method A, change this to: %h/ClaudeHeartbeat/heartbeat.sh
ExecStart=%h/ClaudeHeartbeat/heartbeat.py
Restart=always
RestartSec=60

[Install]
WantedBy=default.target
